# consumers/fasttext/supervisord.conf
[supervisord]
nodaemon=true
logfile_maxbytes=0
loglevel=debug

[program:predictor]
command=/consumers/fasttext/src/main_predictor.py

[program:cache]
command=/consumers/fasttext/src/main_cache.py